# files/etc/uci-defaults/22_crontab
#!/bin/sh
crontab /etc/tube/crontab-input
rm /etc/tube/crontab-input
service cron enable
service tube enable

